# server/ble/sentryusb-ble.service
[Unit]
Description=SentryUSB BLE Peripheral
Documentation=https://github.com/Sentry-Six/Sentry-USB-Rusty
After=bluetooth.target bluetooth.service dbus.service
Requires=bluetooth.target dbus.service
BindsTo=bluetooth.service
ConditionPathExists=/usr/bin/python3
# Defense in depth against the BLE adapter boot race: systemd defaults to
# StartLimitBurst=5 within StartLimitIntervalSec=10s, which can be exhausted
# if all 5 attempts hit the race within 30 seconds. Bumping to 20 attempts in
# 120s gives systemd enough budget to recover before giving up. The Python
# wait_for_adapter() handles 99% of cases on the first attempt; this is
# insurance for unrelated future failure modes (BlueZ crash, OOM, etc.).
StartLimitBurst=20
StartLimitIntervalSec=120

[Service]
Type=simple
# Clean up stale keep-awake locks only (older than 24 h).  Matches the
# BLE_LOCK_MAX_AGE check in run/awake_start.  Unconditionally removing
# the lock here would break mid-session service restarts: if an archive
# cycle is active, awake_stop's GATT-restore flow depends on the lock
# being present.  Wiping a fresh lock leaves the GATT daemon down until
# the next Pi reboot.
ExecStartPre=-/bin/bash -c 'f=/tmp/ble_keep_awake_active; [ -e "$f" ] && [ $(( $(date +%s) - $(stat -c %Y "$f" 2>/dev/null || echo 0) )) -gt 86400 ] && rm -f "$f"; exit 0'
ExecStartPre=/usr/sbin/rfkill unblock bluetooth
ExecStartPre=-/usr/bin/hciconfig hci0 up
# Re-enable controller-level LE advertising. The btmgmt flag persists in
# kernel state and ends up OFF after certain restart paths (service crash,
# manual restart, bond-dir reset). With it OFF, BlueZ's RegisterAdvertisement
# below has no effect — the radio stops broadcasting. Already-bonded phones
# can still reconnect via the bonded address, so the bug is invisible to
# existing users, but any NEW device trying to discover the Pi (first-time
# setup, re-pair after bond wipe, additional family-member phone, future
# iOS client) will silently fail with "Pi not found nearby." The leading
# '-' makes systemd treat a missing btmgmt as non-fatal.
ExecStartPre=-/usr/bin/btmgmt advertising on
ExecStartPre=/bin/bash -c 'for i in $(seq 1 30); do busctl status org.bluez >/dev/null 2>&1 && break; echo "Waiting for org.bluez... ($i/30)"; sleep 1; done'
ExecStart=/usr/bin/python3 /root/bin/sentryusb-ble.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentryusb-ble
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
