# server/ble/sentryusb-telemetry.service
[Unit]
Description=SentryUSB BLE Telemetry Sampler
Documentation=https://github.com/Sentry-Six/Sentry-USB-Rusty
After=bluetooth.target bluetooth.service dbus.service sentryusb.service
Requires=bluetooth.target
# Wait for the main sentryusb daemon so the drive-data DB has had its
# migrate() run — the sampler relies on the v6 schema being in place
# before it starts inserting into telemetry_samples.
After=sentryusb.service
# Won't start if the BLE binaries haven't been installed yet. The
# sampler binary itself, tesla-control (used to talk to the car), and
# the BLE keypair (set during pairing) all need to be on disk before
# we even attempt to run. Listing each separately means
# `systemctl status` says exactly which one is missing rather than
# spawning and crashing at ExecStart.
ConditionPathExists=/root/bin/sentryusb-tesla-telemetry
ConditionPathExists=/root/bin/tesla-control
ConditionPathExists=/root/.ble/key_private.pem
# Match sentryusb-ble.service's burst limits so a flaky boot doesn't
# burn through the systemd restart budget.
StartLimitBurst=20
StartLimitIntervalSec=120

[Service]
Type=simple
# Clean up our own stale lock at startup (only if we own it AND it's
# >24h old) so a crash mid-session doesn't leave the radio held forever.
ExecStartPre=-/bin/bash -c 'f=/tmp/ble_radio_owner; if [ -e "$f" ] && [ "$(head -n1 "$f" 2>/dev/null)" = "telemetry" ] && [ $(( $(date +%s) - $(stat -c %Y "$f" 2>/dev/null || echo 0) )) -gt 86400 ]; then rm -f "$f"; fi; exit 0'
# Source sentryusb.conf so BLE_ENABLED and TESLA_BLE_VIN are visible
# inside the daemon's env. The Rust code re-reads on each iteration,
# but seeding the initial env here makes the systemd status output
# show the active VIN.
#
# We can't use systemd's `EnvironmentFile=` directive: its parser only
# accepts bare `KEY=value` and silently rejects every line starting
# with `export KEY=value`. The user-facing /root/sentryusb.conf is
# bash-style (archiveloop and other legacy shell scripts source it
# directly), so EnvironmentFile= leaves the daemon with an empty env.
#
# Source via bash instead: `set -a` auto-exports every var sourced,
# `source` understands `export` as a no-op when sourcing, and `exec`
# swaps to the binary so systemd's PID tracking still works.
#
# Each file is gated by `[ -f ] && bash -n` so a hand-edited conf with
# a fat-fingered quote can't crash ExecStart and trigger the same
# Restart=always / StartLimitBurst restart-cycling-to-inactive failure
# mode this commit set out to fix. A broken file is skipped (bash -n
# prints the parse error to the journal so the user sees what went
# wrong); the binary still starts and falls back to its own
# parse_file() reader, which is robust to malformed lines.
ExecStart=/bin/bash -c 'set -a; for f in /root/sentryusb.conf /boot/firmware/sentryusb.conf; do [ -f "$f" ] && bash -n "$f" && source "$f"; done; set +a; exec /root/bin/sentryusb-tesla-telemetry'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=sentryusb-telemetry
# Send SIGTERM and give the process up to 8s to release the radio
# lock and restart sentryusb-ble before SIGKILL. The signal handler
# in main.rs does both before returning.
KillSignal=SIGTERM
TimeoutStopSec=8

[Install]
WantedBy=multi-user.target
